# indexer/supervisord.conf
[program:node]
command = node load_environment.js indexer.js
directory = /home/dotcloud/current