# web/supervisord.conf
[program:node]
command = node load_environment.js server.js
directory = /home/dotcloud/current